# contrib/supervisor.ini
[program:mqtt2pushover]
directory = /path/to/mqtt2pushover
command = /path/to/mqtt2pushover.py
user = jpm
environment= MQTT2PUSHOVERCONF="/path/to/mqtt2pushover.conf"
