# etc/supervisor.ini
[program:mqttwarn]
directory = /opt/mqttwarn
command = /opt/mqttwarn/bin/mqttwarn
user = mqttwarn
environment= MQTTWARNINI="/path/to/mqttwarn.ini"
