# etc/supervisor.ini
[program:mqttwarn]
directory = /path/to/mqttwarn
command = /path/to/mqttwarn.py
user = jpm
environment= MQTTWARNINI="/path/to/mqttwarn.ini"
