# etc/supervisor.ini
[program:mqttwarn]
command = /usr/local/bin/mqttwarn
user = mqttwarn
environment = MQTTWARNINI="/path/to/mqttwarn.ini"
